# ansible/roles/borgbackup_config_backup/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup
[Timer]
OnUnitActiveSec=5min
unit=borg-backup.service
[Install]
WantedBy=timers.target
